# claim.ini
[program:claim]
command=/cmd.sh
stdout_logfile=/local/claim.stdout
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/local/claim.stdout
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
autostart=true
autorestart=true
