# consul-template.ini
[program:consul-template]
command=consul-template -config=/config/consul-template.hcl -log-level=info
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
autostart=true
autorestart=true
