# consul.ini
[program:consul]
command=consul agent -config-file=/config/consul.json
autostart=true
autorestart=true

stdout_logfile=/local/consul.stdout
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/local/consul.stderr
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
